# scripts/systemd/speaker-agent.service
[Unit]
Description=Bluetooth speaker agent

[Service]
ExecStart=python ~/bin/speaker-agent.py

[Install]
WantedBy=default.target
